# Tools/gowr's Loud-Volume-Notifier/Related Service Files/mic-alarm.service
[Unit]
Description=Warns when mic volume is too high

[Service]
ExecStart=/home/siyia/.local/bin/mic-alarm/level.out

[Install]
WantedBy=default.target
