# Tools/ChatGPT/Delete Vacant Folders/clean.service
[Unit]
Description=Deletes vacant folders and zip&rar files.

[Service]
Type=oneshot
ExecStart=/bin/delete_vacant_folders

[Install]
WantedBy=default.target
