# MGLRU/mglru_oomd.service
[Unit]
Description=Enable mglru oom killer

[Service]
Type=oneshot
ExecStart=/usr/bin/mglru

[Install]
WantedBy=multi-user.target
