# Zram/zram.service
[Unit]
Description=Enable Zram

[Service]
Type=oneshot
ExecStart=/usr/bin/zram

[Install]
WantedBy=multi-user.target
