# Templates/Timers/timer_daily.timer
[Unit]
Description=A timer for a service that runs everyday at 7:00

[Timer]
OnCalendar=7:00
Persistent=true

[Install]
WantedBy=timers.target
