# Tools/gowr's Loud-Volume-Notifier/Service&Timer/crontab
*/5  * * * * killall level && killall level
