# Tools/gowr's Loud-Volume-Notifier/Service&Timer/mic-alarm.timer
[Unit]
Description=Runs mic alarm service in a timely manner

[Timer]
OnBootSec=1min
OnUnitActiveSec=5sec

[Install]
WantedBy=timers.target
