# Tools/gowr's Loud-Volume-Notifier/Service&Timer/mic-alarm.service
[Unit]
Description=Warns when mic volume is too high

[Service]
Type=oneshot
ExecStart=/home/siyia/.local/bin/mic-alarm/level.out

[Install]
WantedBy=default.target
