# Quirks/Pure-Systemd-Waybar-Launch/login-session.target
[Unit]
Description=Sway Session Target
After=default.target
