# Quirks/Pure-Systemd-Waybar-Launch/waybar.service
[Unit]
Description=Waybar

[Service]
ExecStart=/usr/bin/waybar
Restart=always
RestartSec=1
Environment=PATH=/usr/local/bin:/usr/bin:/bin
